# OpenVMSys.service
[Unit]
Description=OpenVMSys service

[Service]
ExecStart=[Path to your OpenVMSys Executable]

[Install]
WantedBy=multi-user.target